# deployment/ocr-app.service
[Unit]
Description=OCR Application Service
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ocr-app/ocr-app-v1
Environment="PATH=/opt/ocr-app/ocr-app-v1/venv/bin:/usr/bin:/usr/local/bin"
Environment="PYTHONPATH=/opt/ocr-app/ocr-app-v1"
EnvironmentFile=/opt/ocr-app/ocr-app-v1/.env

# Commande de démarrage
ExecStart=/opt/ocr-app/ocr-app-v1/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000

# Configuration mémoire
MemoryMax=8G
MemoryHigh=6G
MemorySwapMax=0
TasksMax=4096

# Configuration des fichiers temporaires
RuntimeDirectory=ocr-app
RuntimeDirectoryMode=0755
ReadWritePaths=/tmp /var/tmp

[Install]
WantedBy=multi-user.target